# systemd/btrfs-scrub@.service.d/ordering.conf
[Unit]
After=btrfs-scrub-resume@%i.service
[Service]
# so can be scheduled after balance
Type=oneshot
IOSchedulingClass=best-effort
