# systemd/backup-hourly-bigdata.service
[Unit]
Description=Hourly backup service of bigdata

[Service]
Type=oneshot
ExecStart=/usr/lib/btrfs-setup/btrfs-backup-bigdata.sh
