# systemd/backup-hourly-bigdata.timer
[Unit]
Description=Hourly backup of bigdata

[Timer]
OnCalendar=hourly
Persistent=true
AccuracySec=1h

[Install]
WantedBy=timers.target
