# systemd/backup-btrbk@daily.timer
[Unit]
Description=Daily backup

[Timer]
OnCalendar=daily
Persistent=true
AccuracySec=1h

[Install]
WantedBy=timers.target
