# systemd/backup-daily-bigdata.timer
[Unit]
Description=Daily backup of bigdata

[Timer]
OnCalendar=daily
Persistent=true
AccuracySec=1h

[Install]
WantedBy=timers.target
