# systemd/backup-btrbk@.service
[Unit]
Description=Backup service on %I

[Service]
Type=oneshot
ExecStart=/usr/bin/btrbk -q run %I
