# systemd/backup-btrbk@.timer
[Unit]
Description=Hourly backup of %I

[Timer]
OnCalendar=hourly
Persistent=true
AccuracySec=1h

[Install]
WantedBy=timers.target
