# systemd/btrfs-balance@.timer
[Unit]
Description=Weekly Btrfs balance on %f

[Timer]
OnCalendar=weekly
AccuracySec=1d
Persistent=true

[Install]
WantedBy=timers.target
